# ansible/roles/nip29_relay/templates/strfry-nip29-sync.service.j2
[Unit]
Description=Strfry NIP-29 negentropy peer sync

[Service]
Type=oneshot
ExecStart={{ nip29_base_dir }}/sync-peers.sh
User=root
