# ansible/roles/nip29_relay/templates/strfry-nip29-sync.timer.j2
[Unit]
Description=Strfry NIP-29 sync timer

[Timer]
OnCalendar=*-*-* *:0/{{ nip29_sync_interval_min }}:00
Persistent=true

[Install]
WantedBy=timers.target
