# ansible/roles/grasp_snapshot/templates/tollgate-grasp-snapshot.timer.j2
[Unit]
Description=Run GRASP snapshot every 5 minutes

[Timer]
OnCalendar=*:0/5
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
