# ansible/roles/grasp_snapshot/templates/tollgate-grasp-snapshot.service.j2
[Unit]
Description=GRASP LMDB + git snapshot for syncthing
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tollgate-grasp-snapshot.sh
