# ansible/roles/backup/templates/tollgate-vps-stats.timer.j2
[Unit]
Description=Tollgate VPS Stats Timer (every 10 seconds)

[Timer]
OnCalendar=*:*:0/10
Persistent=true

[Install]
WantedBy=timers.target
