# ansible/roles/backup/templates/tollgate-vps-stats.service.j2
[Unit]
Description=Tollgate VPS Stats Generator

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gen-vps-stats.py
StandardOutput=journal
StandardError=journal
