# ansible/roles/fips_exit_node/templates/tollgate-fips-exit-advert.timer.j2
[Unit]
Description=Tollgate FIPS exit node route advertisement timer

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ fips_exit_advert_interval_hours }}h
Persistent=true

[Install]
WantedBy=timers.target
